# setup/user-mode-libvirt/files/systemd/user/libvirtd.service
[Unit]
Description=user-mode libvirt daemon (sockets under /run/libvirt-sandbox)
# Created at boot by /etc/tmpfiles.d/libvirt-sandbox.conf.
ConditionPathIsDirectory=/run/libvirt-sandbox

[Service]
Type=notify
# Monolithic daemon on purpose: every driver (qemu, storage, network, secret)
# is served on the one forwarded socket. No --timeout: stay up even with no
# client connected.
ExecStart=/usr/sbin/libvirtd
Restart=on-failure
# Leave qemu guests running across daemon restarts.
KillMode=process

[Install]
WantedBy=default.target
