# config/containers/systemd/forwardproxy.container
# Quadlet for the credential-less forward proxy (rootless / systemd --user).
# Installed by Home Manager to ~/.config/containers/systemd/forwardproxy.container.
#
# Machine-specific office settings live outside git:
#   ~/.config/forwardproxy/forwardproxy.env
#
# %t = $XDG_RUNTIME_DIR, %h = $HOME (systemd specifiers).
[Unit]
Description=Credential-less Kerberos forward proxy (office egress)
Wants=network-online.target forwardproxy-gpg-agent.service
After=network-online.target forwardproxy-gpg-agent.service

[Container]
ContainerName=forwardproxy
# Keep this tag in sync with container/forwardproxy/tag.
Image=localhost/devenv/forwardproxy:0.1.0

# Host networking lets host tools use http://127.0.0.1:3128. Keep the default
# bind localhost-only; widen BIND/PROXY_ACL in forwardproxy.env only where needed.
Network=host

# Forward the dedicated agent's restricted extra socket into the path the
# container's gpg expects. The dedicated keyring should hold only the proxy key.
Volume=%t/forwardproxy/S.gpg-agent.extra:/root/.gnupg/S.gpg-agent:Z

# gpg-encrypted domain password + the public key.
Volume=%h/.config/forwardproxy/secret.gpg:/secrets/secret.gpg:ro,Z
Volume=%h/.config/forwardproxy/pubkey.asc:/secrets/pubkey.asc:ro,Z

# Optional: reuse the host's Kerberos realm/KDC config instead of generating one.
#Volume=/etc/krb5.conf:/etc/krb5.conf:ro

# Required office-specific settings are loaded from an untracked file. The
# systemd drop-in skips startup when the file is absent or the upstream host is
# not resolvable from the current network.
EnvironmentFile=-%h/.config/forwardproxy/forwardproxy.env

# Safe defaults. Override in forwardproxy.env for the office environment.
Environment=BIND=127.0.0.1
Environment=PORT=3128
Environment=PROXY_ACL=127.0.0.0/8
Environment=UPSTREAM_PORT=8080
Environment=PROXY_SPN=HTTP

[Service]
Restart=always
# Allow time for host-side pinentry on (re)start.
TimeoutStartSec=120

[Install]
WantedBy=default.target
