# container/forwardproxy/script/forwardproxy-gpg-agent.service
# Dedicated gpg-agent for the forward proxy (systemd --user).
# Managed copy lives at config/systemd/user/forwardproxy-gpg-agent.service.
# This file is kept as a standalone example for manual installs.
#
# Runs a SEPARATE agent over a SEPARATE keyring that holds only the proxy key,
# and binds a deterministic extra socket the proxy container mounts. This is the
# isolation boundary: the container can ask this agent to decrypt, and that key
# can decrypt only the proxy secret -- your main keys are never reachable.
#
# Adjust the gpg-agent path if not installed via home-manager / nix profile.
[Unit]
Description=Dedicated single-key gpg-agent for forwardproxy
Documentation=man:gpg-agent(1)
ConditionPathExists=%h/.config/forwardproxy/gnupg

[Service]
Type=forking
Environment=FP_GNUPGHOME=%h/.config/forwardproxy/gnupg
# Do not inherit shell/tmux/zellij pinentry routing into this daemon. The proxy
# prompt must use GUI pinentry; terminal pinentry has no service TTY here.
Environment=PINENTRY_USER_DATA=
# Socket lives in the runtime dir (tmpfs, 0700, auto-cleared), not in the
# persistent keyring. systemd creates %t/forwardproxy and removes it on stop.
RuntimeDirectory=forwardproxy
RuntimeDirectoryMode=0700
# GUI pinentry needs a display. With keep-display in gpg-agent.conf the agent
# uses ITS OWN startup display (not the container client's empty one), so these
# must be set when this service starts. Easiest: have your graphical login run
#   systemctl --user import-environment DISPLAY WAYLAND_DISPLAY XAUTHORITY
# (the agent then inherits them from the user manager). Or hardcode below —
# WSLg defaults are DISPLAY=:0 / WAYLAND_DISPLAY=wayland-0.
#Environment=DISPLAY=:0
#Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=%h/.nix-profile/bin/gpg-agent \
    --homedir %h/.config/forwardproxy/gnupg \
    --extra-socket %t/forwardproxy/S.gpg-agent.extra \
    --keep-display \
    --keep-tty \
    --daemon
ExecStop=%h/.nix-profile/bin/gpgconf --homedir %h/.config/forwardproxy/gnupg --kill gpg-agent
Restart=on-failure

[Install]
WantedBy=default.target
