# config/systemd/user/dotfiles-server.service
[Unit]
Description=dotfiles server (devenv build/pull REST API)
Documentation=https://github.com/ngicks/dotfiles
After=network-online.target
Wants=network-online.target

# Skip when moonbit tool chain isn't available. 
ConditionFileIsExecutable=%h/.local/share/moonbit/bin/moon

[Service]
Type=simple
WorkingDirectory=%h/.dotfiles

Environment=MOON_HOME=%h/.local/share/moonbit
Environment=PATH=%h/.local/share/moonbit/bin:%h/.local/share/mise/shims:%h/.local/share/go/bin:%h/.nix-profile/bin:%h/.local/bin:%h/.local/containers/libexec/podman:%h/.local/containers/lib/podman:%h/.local/containers/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

RuntimeDirectory=dotfiles-server
RuntimeDirectoryMode=0700

ExecStartPre=%h/.local/share/moonbit/bin/moon build --target native
ExecStart=%h/.local/share/moonbit/bin/moon run src2 --target native -- server serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
