# config/systemd/user/dotfilesmgr-server.service
[Unit]
Description=dotfilesmgr server (devenv build/pull REST API)
Documentation=https://github.com/ngicks/dotfiles
After=network-online.target
Wants=network-online.target

# Skip when the daemon isn't installed (managed through mise; see
# config/mise/mise.toml).
ConditionFileIsExecutable=%h/.local/share/mise/shims/dotfilesmgr

[Service]
Type=simple
WorkingDirectory=%h/.dotfiles

Environment=PATH=%h/.local/share/mise/shims:%h/.local/share/go/bin:%h/.nix-profile/bin:%h/.local/bin:%h/.local/containers/libexec/podman:%h/.local/containers/lib/podman:%h/.local/containers/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# moonbitlang/async loads libssl with dlopen(3), so it is absent from the
# executable's ELF dependencies and must be discoverable at runtime.
Environment=LD_LIBRARY_PATH=%h/.local/lib/dotfilesmgr-openssl

RuntimeDirectory=dotfilesmgr
RuntimeDirectoryMode=0700

# A login shell raises the soft open-file limit to 10240, but a systemd user
# service keeps systemd's 1024 default, and podman propagates that soft limit
# into the build container. The parallel nix/home-manager build (gcc, the Go
# toolchain, linkers) then exhausts 1024 fds and crashes mid-compile (the
# SIGSEGV compiler ICEs), while the same build from a shell succeeds. Give the
# service the same open-file headroom the manual build gets.
LimitNOFILE=1048576

ExecStart=%h/.local/share/mise/shims/dotfilesmgr server serve
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
