# config/systemd/user/forwardproxy.service.d/10-office-only.conf
# Skip the generated Quadlet service outside the office environment.
[Unit]
ConditionPathExists=%h/.config/forwardproxy/forwardproxy.env

[Service]
ExecCondition=/bin/sh -c '. "$HOME/.config/forwardproxy/forwardproxy.env"; : "${UPSTREAM_PORT:=8080}"; [ -n "$UPSTREAM_HOST" ] && nc -z -w 3 "$UPSTREAM_HOST" "$UPSTREAM_PORT"'
